# systemd/ubuntu18.04/scheduler_motd.timer
[Unit]
Description=Bot Scheduler MOTD

[Timer]
OnCalendar=*-*-* 7:30:00
Persistent=true

[Install]
WantedBy=timers.target
  