# systemd/ubuntu20.04/scheduler_timer.timer
[Unit]
Description=Bot Scheduler timer

[Timer]
OnCalendar=*:0/8 
Persistent=true

[Install]
WantedBy=timers.target
  